# .config/systemd/user/mouse-battery-check.service
[Unit]
Description=Mouse Battery Check

[Service]
Type=oneshot
ExecStart=%h/.bin/mouse-battery-check.sh
