# .config/systemd/user/mouse-battery-check.timer
[Unit]
Description=Check mouse battery every 15 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target
