# .config/systemd/user/wtf-backup.service
[Unit]
Description=WTF Backup
After=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.bin/wtf-backup.sh
