# .config/systemd/user/wtf-backup.timer
[Unit]
Description=Run WTF backup daily at 4 PM

[Timer]
OnCalendar=*-*-* 16:00:00
Persistent=true

[Install]
WantedBy=timers.target
